# deploy/supervisor.conf
; ============================================================
;  Supervisor configuration for my-job Laravel queue workers
;  Place this file at: /etc/supervisor/conf.d/my-job.conf
;  Then run: supervisorctl reread && supervisorctl update
; ============================================================

[program:my-job-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/my-job/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/my-job/storage/logs/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=3600

[program:my-job-scheduler]
; Laravel scheduler — runs every minute, scheduler decides what to execute
process_name=%(program_name)s
command=bash -c "while [ true ]; do php /var/www/my-job/artisan schedule:run --no-interaction >> /dev/null 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/my-job/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
